# debian/pt-sys-menu.service
[Unit]
Description=pi-top OLED system menu

[Service]
Restart=always
ExecStart=/usr/lib/pt-sys-menu/pt-sys-menu

[Install]
WantedBy=local-fs.target
