# debian/pt-sys-oled.postinst
#!/bin/bash -e

systemctl stop hostapd
systemctl disable hostapd
